# service/tweet_itunes/run
!/bin/sh

PATH=/bin:/usr/bin:/sbin:/usr/sbin:/usr/local/bin
export PATH

exec 2>&1
exec setuidgid run ruby /home/run/www/harrier/scripts/twitter/twitter_itunes.rb